# packaging/systemd/replayshield.service
[Unit]
Description=ReplayShield PAM HTTP Auth Server
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/var/lib/replayshield
ExecStartPre=/usr/bin/test -f /dev/shm/replayshield/admin.key
ExecStart=/usr/bin/replayshield serve
Restart=no